# Pi/systemd/raptorhab-ground.service
[Unit]
Description=RaptorHab Ground Station
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/raptorhab
ExecStart=/usr/bin/python3 -m ground.main
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment="RAPTORHAB_GND_FREQUENCY=915.0"
Environment="RAPTORHAB_GND_DATA_PATH=/home/pi/raptorhab-ground/data"
Environment="RAPTORHAB_GND_IMAGE_PATH=/home/pi/raptorhab-ground/images"
Environment="RAPTORHAB_GND_LOG_PATH=/home/pi/raptorhab-ground/logs"
Environment="RAPTORHAB_GND_WEB_PORT=5000"

# Resource limits
MemoryMax=512M
CPUQuota=80%

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/raptorhab-ground

[Install]
WantedBy=multi-user.target
