# systemd/raptorhab-airborne.service
[Unit]
Description=RaptorHab Airborne Payload
Documentation=https://github.com/raptorhab/raptorhab
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/raptorhab
ExecStart=/usr/bin/python3 -m airborne.main
Restart=always
RestartSec=10
TimeoutStopSec=30

# Resource limits
MemoryMax=256M
CPUQuota=80%

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/pi/raptorhab

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=raptorhab-airborne

[Install]
WantedBy=multi-user.target
